# root/usr/lib/systemd/system/univention-container-mode-firstboot.service.d/override.conf
[Unit]
OnFailure=univention-container-mode-firstboot-on-failure.service
