# root/usr/lib/systemd/system/ssh.service.d/override.conf
[Service]
# /run are tmpfs, we have to create a path
ExecStartPre=
ExecStartPre=/bin/sh -c '[ ! -d /run/sshd ] || /bin/true && /bin/mkdir --parents /run/sshd'
ExecStartPre=/usr/sbin/sshd -t
