# root/usr/lib/systemd/system/networking.service.d/override.conf
[Service]
ExecStartPre=
ExecStartPre=/bin/true
ExecStart=
ExecStart=/bin/true
ExecStop=
ExecStop=/bin/true
