# root/usr/lib/systemd/system/univention-container-mode-recreate-on-failure.service.d/override.conf
[Unit]
ConditionPathExists=
ConditionPathExists=/var/backups/univention-container-mode/restore
ConditionPathExists=/var/backups/univention-container-mode
ConditionPathExists=/dev/shm/univention-container-mode.env
ConditionPathExists=/usr/lib/univention-container-mode
ConditionPathExists=!/var/univention-join/joined
ConditionPathExists=!/var/univention-join/status
ConditionPathExists=!/etc/univention/base.conf

[Service]
ExecStartPre=

[Install]
WantedBy=
