# root/usr/lib/systemd/system/univention-container-mode-recreate.service.d/override.conf
[Unit]
OnFailure=univention-container-mode-recreate-on-failure.service
