# root/usr/lib/systemd/system/ntpsec.service.d/override.conf
[Service]
ExecStartPre=-/bin/bash -c '[[ -d /var/log/ntpsec ]] || mkdir /var/log/ntpsec'
ExecStartPre=-/usr/bin/sed -i 's/^NTPD_OPTS="$NTPD_OPTS -u ntpsec:ntpsec"/NTPD_OPTS="$NTPD_OPTS"/g' /usr/libexec/ntpsec/ntp-systemd-wrapper
