# root/usr/lib/systemd/system/named.service.d/override.conf
[Service]
# force cache flush on startup ( named managed-keys.bind.jnl: open: permission denied afer installing samba )
ExecStartPre=-/bin/rm --force --verbose \
  /var/cache/bind/managed-keys.bind.jnl \
  /var/cache/bind/managed-keys.bind
# set working directory for both bind name server ( named.conf.proxy or named.conf.samba4 )
WorkingDirectory=/var/cache/bind
