# root/usr/lib/systemd/system/avahi-daemon.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/sbin/avahi-daemon --syslog --no-chroot
ExecReload=
ExecReload=/usr/sbin/avahi-daemon --reload
