# root/usr/lib/systemd/system/slapd.service.d/override.conf
[Service]
# /var/run are tmpfs, we have to create a path
ExecStartPre=/bin/sh -c '[ ! -d /var/run/slapd ] || /bin/true && /bin/mkdir --parents /var/run/slapd'
