# root/usr/lib/systemd/system/univention-container-mode-firstboot-on-failure.service.d/override.conf
[Unit]
Description=
Description=Univention container mode firstboot ( on failure )
ConditionPathExists=
ConditionPathExists=/dev/shm/univention-container-mode.env
ConditionPathExists=/usr/lib/univention-container-mode
ConditionPathExists=!/var/univention-join/joined
ConditionPathExists=!/var/univention-join/status
ConditionPathExists=!/etc/machine.secret

[Install]
WantedBy=
