# root/usr/lib/systemd/system/auth-rpcgss-module.service.d/override.conf
[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=
ExecStart=-/sbin/modprobe -q auth_rpcgss
ExecStart=/bin/bash -c '/bin/lsmod | /bin/egrep --quiet -- auth_rpcgss'