# root/usr/lib/systemd/system/univention-saml.service.d/override.conf
[Service]
# /var/run are tmpfs, we have to create a path
ExecStartPre=/bin/sh -c '[ ! -d /var/run/univention-saml ] || /bin/true && ( /bin/mkdir --parents /var/run/univention-saml && /bin/chown samlcgi.memcache /var/run/univention-saml )'
