# root/usr/lib/systemd/system/univention-container-mode-pre-installed-role.service.d/override.conf
[Unit]
Description=
Description=Univention container mode firstboot with pre installed role as %%ROLE%%
ConditionPathExists=
ConditionPathExists=/dev/shm/univention-container-mode.env
ConditionPathExists=/usr/lib/univention-container-mode
ConditionPathExists=!/var/univention-join/joined
ConditionPathExists=!/var/univention-join/status
ConditionPathExists=/etc/univention/base.conf
OnFailure=univention-container-mode-pre-installed-role-on-failure.service

[Service]
ExecStartPre=
ExecStartPre=-/bin/bash -c 'echo "role=%%ROLE%%" >> /dev/shm/univention-container-mode.env'

[Install]
WantedBy=
