# root/usr/lib/systemd/system/nscd.service.d/override.conf
[Unit]
Wants=slapd.service
After=slapd.service

[Service]
# /var/run are tmpfs, we have to create a path
ExecStartPre=/bin/sh -c '[ ! -d /run/nscd ] || /bin/true && /bin/mkdir --parents /run/nscd'
