# root/usr/lib/systemd/system/nagios-nrpe-server.service.d/override.conf
[Service]
User=root
Group=root
# /var/run are tmpfs, we have to create a path
ExecStartPre=/bin/sh -c '[ ! -d /var/run/nagios ] || /bin/true && /bin/mkdir --parents /var/run/nagios'
